# res/etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                  ; run supervisord as a daemon

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:dalton]
command = python /opt/dalton/run.py
directory = /opt/dalton/

[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true
;stderr_logfile = /var/log/supervisord/dalton-stderr.log
;stdout_logfile = /var/log/supervisord/dalton-stdout.log
